# backup_tend.service
[Unit]
Description=Backup tend databases

[Service]
Type=oneshot
ExecStart=/home/isaac/code/tend/backup_tend
